# files/etc/uci-defaults/luci-xlnetacc
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@xlnetacc[-1]
	add ucitrack xlnetacc
	set ucitrack.@xlnetacc[-1].init=xlnetacc
	commit ucitrack
EOF

general=$(uci -q get xlnetacc.@general[-1])
if [ -z "$general" ]; then
	uci -q add xlnetacc general
fi
if [ "$general"x != "general"x ]; then
	uci -q batch <<-EOF >/dev/null
		rename xlnetacc.@general[-1]="general"
		commit xlnetacc
	EOF
fi

# 配置项列表
config_options="enabled:0 down_acc:0 up_acc:0 logging:1 verbose:0 network:wan keepalive:10 relogin:0 account: password: base_url:https://openrouter.ai/api/v1 api_key: model:google/gemini-2.0-flash-exp:free"

# 检查并设置配置项
need_commit=0
uci_batch=""

for option in $config_options; do
	key="${option%%:*}"
	value="${option#*:}"
	
	if ! uci -q get xlnetacc.@general[-1].$key >/dev/null 2>&1; then
		uci_batch="${uci_batch}	set xlnetacc.@general[-1].$key='$value'
"
		need_commit=1
	fi
done

if [ $need_commit -eq 1 ]; then
	uci -q batch <<-EOF >/dev/null
$uci_batch	commit xlnetacc
	EOF
fi

rm -rf /tmp/luci-indexcache /tmp/luci-modulecache
exit 0
